# rootfs/systemd-services/etc/systemd/system/systemd-suspend.service.d/override.conf
[Service]
ExecStartPost=systemctl restart bluetooth-wakeup.service
ExecStartPre=systemctl restart bluetooth-wakeup.service
